# crontab
# min hour day month week  cmd
  *   */4  *  *      *    /root/app/send2kindle/run.sh 
