# crontab
# min hour day month week  cmd
  *   */4  *  *      *   cd /app/send2kindle/src && python3 main.py && rm *.txt 
